# launchd/com.ccfleet.ttyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE — do not copy this file directly to /Library/LaunchDaemons/.
  Run  sudo bash scripts/install-launchd.sh  to generate the real plist
  with your username and paths substituted in.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ccfleet.ttyd</string>

  <key>UserName</key>
  <string>__USER__</string>

  <key>ProgramArguments</key>
  <array>
    <string>__INSTALL_DIR__/bin/start-ttyd.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__INSTALL_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/ccfleet/ttyd.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/ccfleet/ttyd-error.log</string>

  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
